# templates/alloy-systemd.service.j2
# ============================================================================
# Alloy systemd unit file
# ============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY!       ⚠️ 
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
#
# Template path: roles/monitoring/templates/alloy-systemd.service.j2
# Deploy path: /etc/systemd/system/alloy.service
# =============================================================================
#
# This systemd service file configures Alloy to run as a system service,
# providing metrics collection and forwarding for observability pipelines.
#
# Features:
#   - Runs Alloy as a background service
#   - Uses dedicated user and group
#
# Configuration:
#   - Config file: {{ alloy_config_dir }}/alloy.yaml
#   - Data directory: {{ alloy_data_dir }}
#   - User: {{ alloy_user }}
#   - Working directory: {{ alloy_data_dir }}
#   - Binary location: /usr/local/bin/alloy
#   - Logs: sudo journalctl -u alloy -f
#
# Usage:
#   sudo systemctl enable alloy    # Enable on boot
#   sudo systemctl start alloy     # Start service
#   sudo systemctl status alloy    # Check status
#   sudo systemctl stop alloy      # Stop service
#   sudo journalctl -u alloy       # View logs
#   sudo journalctl -u alloy -f    # View and follow logs
#
# Security Notes:
#   - Alloy runs with limited privileges
#
# Dependencies:
#   - Requires network.target
#
# References:
#   - https://github.com/grafana/alloy
#
# Generated by Ansible on {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
#
# ------------------------------------------------------------------------------

[Unit]
Description=Alloy Service
Documentation=https://github.com/grafana/alloy

# Service Dependencies
After=network.target
Requires=network.target

[Service]
# The folder to run the ExecStart command from
WorkingDirectory={{ alloy_data_dir }}

ExecStart=/usr/local/bin/alloy run {{ alloy_config_file }} \
                               --server.http.listen-addr "{{ monitoring_alloy.listen_address }}:{{ monitoring_alloy.http_port }}" \
                               --server.http.ui-path-prefix "/{{ monitoring_alloy.path_prefix }}" \
                               --storage.path "{{ alloy_data_dir }}"

# Process management
Type=simple
Restart=on-failure
RestartSec=10

# Directory creation and permissions
User={{ alloy_user }}
Group={{ alloy_group }}
RuntimeDirectory={{ alloy_user }}
RuntimeDirectoryMode=0750
ReadOnlyPaths={{ alloy_config_dir }}
ReadWritePaths={{ alloy_data_dir }}

# Logging to systemd
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ alloy_user }}

# Hardening measures
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateDevices=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
