# templates/grafana-systemd.service.j2
# ============================================================================
# Grafana Systemd Service Unit
# ============================================================================
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY!       ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
#
# Template path: roles/monitoring/templates/grafana-systemd.service.j2
# Deploy path: /etc/systemd/system/{{ grafana_systemd_name }}.service
# Last deployed: {{ ansible_date_time.iso8601 }}
# =============================================================================
#
# This systemd service file configures Grafana to run as a system service,
# providing dashboards and visualization for observability data.
#
# Features:
#   - Runs Grafana as a background service
#   - Uses dedicated user and group
#   - Hardened with systemd security options
#
# Configuration:
#   - Config file: {{ grafana_config_file }}
#   - Data directory: {{ grafana_data_dir }}
#   - User: {{ grafana_user }}
#   - Working directory: {{ grafana_data_dir }}
#   - Binary location: {{ grafana_binary_path }}
#   - Logs: sudo journalctl -u {{ grafana_systemd_name }} -f
#
# Usage:
#   sudo systemctl enable {{ grafana_systemd_name }}    # Enable on boot
#   sudo systemctl start {{ grafana_systemd_name }}     # Start service
#   sudo systemctl status {{ grafana_systemd_name }}    # Check status
#   sudo systemctl stop {{ grafana_systemd_name }}      # Stop service
#   sudo journalctl -u {{ grafana_systemd_name }}       # View logs
#   sudo journalctl -u {{ grafana_systemd_name }} -f    # View and follow logs
#
# Security Notes:
#   - Grafana runs with limited privileges
#
# Dependencies:
#   - Requires network.target
#
# References:
#   - https://grafana.com/docs/grafana/latest/setup-grafana/configure-grafana/
#
# Generated by Ansible on {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# ------------------------------------------------------------------------------

[Unit]
Description=Grafana Service
Documentation=https://grafana.com/docs/grafana/latest/

# Service Dependencies
After=network.target postgres.service
Requires=network.target

[Service]
# The folder to run the ExecStart command from
WorkingDirectory={{ grafana_home_dir }}
Environment="GF_HOME={{ grafana_home_dir }}"

ExecStart={{ grafana_binary_path }} server \
                                    --config={{ grafana_config_file }} \
                                    --pidfile=/run/{{ grafana_user }}/grafana-server.pid \
                                    cfg:default.paths.data={{ grafana_data_dir }} \
                                    cfg:default.paths.plugins={{ monitoring_grafana.plugin_dir }} \
                                    cfg:default.paths.provisioning={{ monitoring_grafana.provisioning_dir }}

Type=simple
Restart=on-failure
RestartSec=10

# Directory creation and permissions
User={{ grafana_user }}
Group={{ grafana_group }}
RuntimeDirectory={{ grafana_user }}
RuntimeDirectoryMode=0750
ReadWritePaths={{ grafana_home_dir }} {{ grafana_config_dir }}

# Logging to systemd
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ grafana_user }}

# Hardening measures
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateDevices=true
#MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET
RestrictRealtime=true

[Install]
WantedBy=multi-user.target
