# templates/loki-systemd.service.j2

# ============================================================================ 
# Loki Systemd Service Unit
# ============================================================================ 
# ⚠️ WARNING: This file is managed by Ansible - DO NOT EDIT MANUALLY!       ⚠️
# Any manual changes will be overwritten on the next Ansible run.
# Edit the template file or Ansible variables instead.
#
# file_path: roles/monitoring/templates/loki-systemd.service.j2
# ============================================================================ 
#
# This systemd service file configures Loki to run as a system service,
# providing log aggregation and querying for distributed systems.
#
# Features:
#   - Runs Loki as a dedicated system user
#   - Uses a minimal, hardened systemd configuration
#
# Configuration:
#   - Config file: {{ loki_config_dir }}/loki-config.yaml
#   - Data directory: {{ loki_data_dir }}
#   - User: {{ loki_user }}
#   - Working directory: {{ loki_data_dir }}
#   - Binary location: /usr/local/bin/loki-linux-amd64
#   - Logs: sudo journalctl -u {{ loki_systemd_name }} -f
#
# Usage:
#   sudo systemctl enable {{ loki_systemd_name }}    # Enable on boot
#   sudo systemctl start {{ loki_systemd_name }}     # Start service
#   sudo systemctl status {{ loki_systemd_name }}    # Check status
#   sudo systemctl stop {{ loki_systemd_name }}      # Stop service
#   sudo journalctl -u {{ loki_systemd_name }}       # View logs
#   sudo journalctl -u {{ loki_systemd_name }} -f    # View and follow logs
#
# Security Notes:
#   - Loki runs with limited privileges
#
# Dependencies:
#   - Requires network.target
#
# References:
#   - https://github.com/grafana/loki
#
# Generated by Ansible on {{ ansible_date_time.iso8601 }}
# Host: {{ ansible_fqdn }}
# ------------------------------------------------------------------------------


# ===============================[Unit]=======================================
[Unit]
Description=Grafana Loki Log Aggregation Service
Documentation=https://github.com/grafana/loki
After=network.target
Requires=network.target


# =============================[Service]======================================
[Service]
# The folder to run the ExecStart command from
WorkingDirectory={{ loki_data_dir }}

# Start Loki with the specified config file
ExecStart={{ loki_binary }} --config.file={{ loki_config_file }}

# Process management
Type=simple
Restart=on-failure
RestartSec=10

# Directory creation and permissions
User={{ loki_user }}
Group={{ loki_group }}
RuntimeDirectory={{ loki_user }}
RuntimeDirectoryMode=0750
ReadOnlyPaths={{ loki_config_dir }}
ReadWritePaths={{ loki_data_dir }}

# Logging to systemd
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ loki_systemd_name }}

# Hardening measures
PrivateTmp=true
ProtectSystem=full
ProtectHome=true
NoNewPrivileges=true
PrivateDevices=true
MemoryDenyWriteExecute=true
SystemCallArchitectures=native
ProtectKernelModules=true
ProtectControlGroups=true
RestrictAddressFamilies=AF_INET
RestrictRealtime=true


# =============================[Install]======================================
[Install]
WantedBy=multi-user.target
